# crontab
5 * * * * zsh ~/github/dotfiles/scripts/check_homebrew_updates.sh
